# wearable-pin/pi/services/camera.service
[Unit]
Description=Wearable Pin Camera Capture Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/wearable-pin/pi
ExecStart=/usr/bin/python3 /home/pi/wearable-pin/pi/capture_image.py --continuous
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
